# extract/crontab
0 */1 * * * source ~/.virtualenvs/base/bin/activate; python ~/gitworkspace/datawarehousing/bgse-dw-final-project/covid_to_mongo.py